# supervisord.conf
[program:node]
command = ./node_modules/.bin/coffee dev/game.coffee
directory = /home/dotcloud/current
